# spoolup.service
[Unit]
Description=SpoolUp - YouTube Streamer for 3D Prints
After=network-online.target moonraker.service
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/usr/data/spoolup
ExecStart=/usr/bin/python3 /usr/data/spoolup/spoolup.py -c /usr/data/spoolup/config.json
Restart=always
RestartSec=10
StandardOutput=append:/var/log/spoolup.log
StandardError=append:/var/log/spoolup.log

[Install]
WantedBy=multi-user.target
